# com.feishu.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.feishu.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shaoshuaima/.nvm/versions/node/v24.13.0/bin/node</string>
        <string>/Users/shaoshuaima/workspace/cursor/feishu-agent/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shaoshuaima/workspace/cursor/feishu-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/Users/shaoshuaima/.nvm/versions/node/v24.13.0/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/shaoshuaima/logs/feishu-agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shaoshuaima/logs/feishu-agent-error.log</string>
</dict>
</plist>
